# scripts/maize.cron
0 3 * * * /bin/sh /var/www/html/MaizeWebApp/scripts/maize_app_maintenance_script.sh > /var/www/html/MaizeWebApp/data/cron.out 2> /var/www/html/MaizeWebApp/data/cron.err 

